# scripts/livox_forwarder.service
[Unit]
Description=Livox UDP forwarder (points + IMU to Pi and localhost)
After=network.target

[Service]
Environment="PI_IP=192.168.2.2"
Environment="PI_PORT=5600"
ExecStart=/home/capra/projects/photorealistic-mapping/scripts/livox_forwarder/livox_forwarder ${PI_IP} ${PI_PORT}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
