# scripts/livox_forwarder.service
[Unit]
Description=Livox point data forwarder (replicates lidar UDP to Pi and localhost)
After=network.target

[Service]
Environment="FORWARDER_SRC=/home/capra/projects/photorealistic-mapping/livox_forwarder"
ExecStartPre=/bin/bash -c "g++ -O2 -o /home/capra/projects/photorealistic-mapping/livox_forwarder/livox_forwarder /home/capra/projects/photorealistic-mapping/livox_forwarder/livox_forwarder.cpp"
ExecStart=/home/capra/projects/photorealistic-mapping/livox_forwarder/livox_forwarder 192.168.1.50 5600 56301 56302
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
